# trmnl-screenshot/rootfs/etc/services.d/trmnl-screenshot/finish
#!/bin/sh
# S6 finish script - called when service terminates

exit 0
