# trmnl-screenshot/rootfs/etc/services.d/trmnl-screenshot/run
#!/command/execlineb -P
# S6-overlay service script for TRMNL Screenshot addon

exec -c node /app/src/index.js
